# crypto-alert-bot/crypto-alert-bot.service
[Unit]
Description=Crypto Alert Bot (Production)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ageek-work/crypto-alert-bot
ExecStart=/usr/bin/python3 src/saas_main.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
